# deploy/lilus.service
[Unit]
Description=LILUS — Sistema de gestion de ventas
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__APP_DIR__
Environment=NODE_ENV=production
Environment=PORT=3000
EnvironmentFile=__APP_DIR__/.env
ExecStart=/usr/bin/node node_modules/next/dist/bin/next start -p 3000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
